# files/etcd.conf
[program:etcd]
command=etcd
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
directory=/files
stdout_logfile=/var/log/supervisor/%(program_name)s/stdout.log
stderr_logfile=/var/log/supervisor/%(program_name)s/stderr.log
